# devices/ssc30kq_lte/general/overlay/etc/systemd/system/lte-setup.service
[Unit]
Description=Huawei E3372h LTE Modem Setup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lte-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
